# bin/cron.daily
#!/bin/bash

buildpath=/srv/buildd

#echo "*** update local repositories with current source ***"

reprepro -b $buildpath/repo/ update precise
